# tools/systemd/telegram-gateway.service
[Unit]
Description=Telegram Gateway - Shared service for Telegram message ingestion
Documentation=https://github.com/yourusername/TradingSystem
After=network-online.target
Wants=network-online.target

# Ensure TimescaleDB is available (if running as container)
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=marce
Group=marce
WorkingDirectory=/home/marce/Projetos/TradingSystem/apps/telegram-gateway

# Environment variables (loaded from .env)
EnvironmentFile=/home/marce/Projetos/TradingSystem/.env

# Start command
ExecStart=/usr/bin/node src/index.js

# Restart policy
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitInterval=300

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/marce/Projetos/TradingSystem/apps/telegram-gateway/data
ReadWritePaths=/home/marce/Projetos/TradingSystem/apps/telegram-gateway/logs

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-gateway

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
